# dist/node.service
[Unit]
Description=NodeJS, Express server

[Service]
Type=simple
ExecStart=/home/ec2-user/.local/bin/node /home/ec2-user/blog/site.js
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
